# supervisord/comfyui.conf
[program:comfyui]
command=python /app/Comfyui/main.py --lowvram --listen
priority=2
autostart=1
autorestart=true
startsecs=0
environment=