# supervisord/cron_script.conf
[program:cron_script]
command=bash /app/cron_script.sh
priority=2
autostart=1
autorestart=1
startsecs=0