# supervisord/dockerd.conf
[program:dockerd]
command=/usr/bin/dockerd
priority=1
autostart=false
autorestart=true
startsecs=0