# supervisord/litellm.conf
[program:litellm]
command=litellm
priority=2
autostart=1
autorestart=true
startsecs=0
environment=