# supervisord/start_it.conf
[program:start_it]
command=bash /app/start_script.sh
priority=2
autostart=1
autorestart=false
startsecs=0